# jingdong/crontest.cron
38 22 * * * sh /Users/doulihang/work/project/weeds_study/jingdong/run.sh
